# contrib/init/groincoind.conf
description "Groincoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env GROINCOIND_BIN="/usr/bin/groincoind"
env GROINCOIND_USER="groincoin"
env GROINCOIND_GROUP="groincoin"
env GROINCOIND_PIDDIR="/var/run/groincoind"
# upstart can't handle variables constructed with other variables
env GROINCOIND_PIDFILE="/var/run/groincoind/groincoind.pid"
env GROINCOIND_CONFIGFILE="/etc/groincoin/groincoin.conf"
env GROINCOIND_DATADIR="/var/lib/groincoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # groincoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$GROINCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run groincoind."
        echo "The setting must appear in $GROINCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Groincoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$GROINCOIND_PIDDIR"
    chmod 0755 "$GROINCOIND_PIDDIR"
    chown $GROINCOIND_USER:$GROINCOIND_GROUP "$GROINCOIND_PIDDIR"
    chown $GROINCOIND_USER:$GROINCOIND_GROUP "$GROINCOIND_CONFIGFILE"
    chmod 0660 "$GROINCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$GROINCOIND_PIDFILE" \
    --chuid $GROINCOIND_USER:$GROINCOIND_GROUP \
    --exec "$GROINCOIND_BIN" \
    -- \
    -pid="$GROINCOIND_PIDFILE" \
    -conf="$GROINCOIND_CONFIGFILE" \
    -datadir="$GROINCOIND_DATADIR" \
    -disablewallet \
    -daemon

